# puppet-files/logstash_upstart.conf
description "logstash"
start on runlevel [2]
stop on runlevel [016]

exec /usr/local/logstash/bin/logstash -f /etc/logstash/logstash.conf -l /var/log/logstash.log --debug

respawn
